# qemu-pr-helper.socket
[Unit]
Description=Persistent Reservation Daemon for QEMU

[Socket]
ListenStream=/run/qemu-pr-helper.sock

[Install]
WantedBy=multi-user.target
